# recipes-utils/archconfig/files/archconfig.service
[Unit]
Description=Package Manager Arch Configuration

[Service]
Type=oneshot
ExecStart=/usr/bin/archconfig.sh

[Install]
WantedBy=multi-user.target
